# templates/launchd/digest-inbox.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>{LABEL_PREFIX}.digest-inbox</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>{HOME}/.claude/skills/digest/scripts/inbox_poll.py</string>
    <string>--once</string>
  </array>
  <key>StartInterval</key><integer>{INTERVAL_SEC}</integer>
  <key>StandardOutPath</key><string>{HOME}/.claude/skills/digest/reports/inbox.log</string>
  <key>StandardErrorPath</key><string>{HOME}/.claude/skills/digest/reports/inbox.err</string>
  <key>RunAtLoad</key><true/>
</dict>
</plist>
